# homechat/rootfs/etc/services.d/avahi-daemon/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Finish script for Avahi daemon
# ==============================================================================

bashio::log.info "Avahi daemon stopped"