# homechat/rootfs/etc/services.d/avahi-daemon/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Avahi daemon for mDNS/DNS-SD service discovery
# ==============================================================================

bashio::log.info "Starting Avahi daemon for mDNS service discovery..."

# Create runtime directory
mkdir -p /var/run/avahi-daemon

# Start Avahi daemon
exec avahi-daemon --no-drop-root --no-chroot --no-rlimits