# homechat/rootfs/etc/services.d/dbus/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# D-Bus system daemon finish script
# ==============================================================================

bashio::log.info "D-Bus daemon stopped"