# homechat/rootfs/etc/services.d/dbus/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start D-Bus system daemon (or use host D-Bus if available)
# ==============================================================================

# Check if host D-Bus is already available (host_dbus: true in config.yaml)
if [ -S /var/run/dbus/system_bus_socket ]; then
    bashio::log.info "Using host D-Bus (socket already available)"
    # Keep service running but don't start our own daemon
    exec sleep infinity
fi

bashio::log.info "Starting D-Bus system daemon..."

# Create runtime directories
mkdir -p /var/run/dbus
mkdir -p /var/lib/dbus

# Initialize D-Bus machine ID if it doesn't exist
if [ ! -f /var/lib/dbus/machine-id ]; then
    dbus-uuidgen > /var/lib/dbus/machine-id
fi

# Start D-Bus system daemon
exec dbus-daemon --system --nofork --nopidfile