# homechat/rootfs/etc/services.d/homechat/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start HomeChat Rails application
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare site_name
declare allow_signups
declare port
declare access_mode
declare ssl
declare certfile
declare keyfile
declare log_level
declare enable_integrations
declare auto_create_api_token
declare home_assistant_integration
declare network_range

bashio::log.info "Starting HomeChat..."

# Get configuration options
site_name=$(bashio::config 'site_name')
allow_signups=$(bashio::config 'allow_signups')
port=$(bashio::config 'port')
access_mode=$(bashio::config 'access_mode')
ssl=$(bashio::config 'ssl')
certfile=$(bashio::config 'certfile')
keyfile=$(bashio::config 'keyfile')
log_level=$(bashio::config 'log_level')
enable_integrations=$(bashio::config 'enable_integrations')
auto_create_api_token=$(bashio::config 'auto_create_api_token')
home_assistant_integration=$(bashio::config 'home_assistant_integration')
network_range=$(bashio::config 'network_range')

# Set environment variables
export RAILS_ENV=production
export RAILS_SERVE_STATIC_FILES=true
export RAILS_LOG_TO_STDOUT=true
export RAILS_RELATIVE_URL_ROOT=""
export HOME_ASSISTANT_ADDON=true

# SSL Configuration based on access mode
# - ingress: HA handles SSL termination, we assume SSL from proxy headers
# - direct_ssl: Terminate SSL at Rails/Puma layer with certificates
# - direct_http: No SSL, plain HTTP (for local network only)
case "${access_mode}" in
    "ingress")
        bashio::log.info "Access mode: Home Assistant Ingress (SSL handled by HA)"
        export RAILS_ASSUME_SSL=true
        export RAILS_FORCE_SSL=false
        ssl_enabled=false
        ;;
    "direct_ssl")
        bashio::log.info "Access mode: Direct SSL (certificates at /ssl)"
        export RAILS_ASSUME_SSL=false
        export RAILS_FORCE_SSL=false
        ssl_enabled=true
        # Validate certificate files exist
        if [[ ! -f "/ssl/${certfile}" ]]; then
            bashio::log.error "SSL certificate not found: /ssl/${certfile}"
            bashio::exit.nok "SSL certificate file missing"
        fi
        if [[ ! -f "/ssl/${keyfile}" ]]; then
            bashio::log.error "SSL key not found: /ssl/${keyfile}"
            bashio::exit.nok "SSL key file missing"
        fi
        bashio::log.info "Using SSL certificates: /ssl/${certfile}, /ssl/${keyfile}"
        ;;
    "direct_http")
        bashio::log.info "Access mode: Direct HTTP (no SSL - local network only)"
        export RAILS_ASSUME_SSL=false
        export RAILS_FORCE_SSL=false
        ssl_enabled=false
        ;;
    *)
        bashio::log.warning "Unknown access mode: ${access_mode}, defaulting to ingress"
        export RAILS_ASSUME_SSL=true
        export RAILS_FORCE_SSL=false
        ssl_enabled=false
        ;;
esac
# Generate or load persistent secret key
if [[ ! -f "/data/secret_key_base" ]]; then
    echo "$(cat /proc/sys/kernel/random/uuid | sed 's/-//g')$(cat /proc/sys/kernel/random/uuid | sed 's/-//g')" > /data/secret_key_base
fi
export SECRET_KEY_BASE=$(cat /data/secret_key_base)
export PORT="${port}"
export SITE_NAME="${site_name:-HomeChat}"
export ALLOW_SIGNUPS="${allow_signups:-true}"
export ENABLE_INTEGRATIONS="${enable_integrations:-true}"
export AUTO_CREATE_API_TOKEN="${auto_create_api_token:-true}"
export HOME_ASSISTANT_INTEGRATION="${home_assistant_integration:-true}"
export NETWORK_RANGE="${network_range:-192.168.0.0/16}"

# Set Home Assistant specific environment variables
if bashio::supervisor.ping; then
    export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN}"
    export HOMEASSISTANT_URL="http://supervisor/core"
fi

# Set log level for Rails
case "${log_level}" in
    "debug")
        export RAILS_LOG_LEVEL=debug
        ;;
    "info")
        export RAILS_LOG_LEVEL=info
        ;;
    "warning")
        export RAILS_LOG_LEVEL=warn
        ;;
    "error")
        export RAILS_LOG_LEVEL=error
        ;;
    *)
        export RAILS_LOG_LEVEL=info
        ;;
esac

bashio::log.info "Site Name: ${SITE_NAME}"
bashio::log.info "Allow Signups: ${ALLOW_SIGNUPS}"
bashio::log.info "Port: ${PORT}"
bashio::log.info "SSL: ${ssl}"
bashio::log.info "Log Level: ${log_level}"

# Ensure database exists
cd /app || bashio::exit.nok "Could not change to app directory"

# Setup persistent storage directory for Active Storage uploads
if [[ ! -d "/data/storage" ]]; then
    bashio::log.info "Creating persistent storage directory..."
    mkdir -p /data/storage
fi

# Symlink app storage to persistent data directory
if [[ -d "/app/storage" && ! -L "/app/storage" ]]; then
    # Move any existing files to persistent storage
    if [[ "$(ls -A /app/storage 2>/dev/null)" ]]; then
        bashio::log.info "Moving existing storage files to persistent location..."
        cp -rn /app/storage/* /data/storage/ 2>/dev/null || true
    fi
    rm -rf /app/storage
fi
if [[ ! -L "/app/storage" ]]; then
    ln -sf /data/storage /app/storage
    bashio::log.info "Linked /app/storage -> /data/storage"
fi

# Setup database if needed
if [[ ! -f "/data/production.sqlite3" ]]; then
    bashio::log.info "Setting up database..."
    export DATABASE_URL="sqlite3:///data/production.sqlite3"
    ./bin/rails db:prepare
else
    bashio::log.info "Database exists, running migrations..."
    export DATABASE_URL="sqlite3:///data/production.sqlite3"
    ./bin/rails db:migrate
fi

# Always run seeds to ensure admin user is created if needed
bashio::log.info "Running database seeds..."
./bin/rails db:seed

# Display admin credentials if they exist
if [[ -f "/data/admin_credentials.json" ]]; then
    bashio::log.info "ðŸ” Admin credentials available!"

    # Parse credentials from JSON file using bashio (Home Assistant's JSON parser)
    admin_username=$(bashio::jq /data/admin_credentials.json '.username')
    admin_password=$(bashio::jq /data/admin_credentials.json '.password')
    created_at=$(bashio::jq /data/admin_credentials.json '.created_at')

    bashio::log.info "â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”"
    bashio::log.info "ðŸš€ HOMECHAT ADMIN ACCESS"
    bashio::log.info "â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”"
    bashio::log.info "ðŸ‘¤ Username: ${admin_username}"
    bashio::log.info "ðŸ”‘ Password: ${admin_password}"
    bashio::log.info "ðŸ•’ Created: ${created_at}"
    bashio::log.info "â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”"
    bashio::log.info "âš ï¸  Please change this password after first login!"
    bashio::log.info "ðŸŒ Access HomeChat at: http://homeassistant.local:${port}"
    bashio::log.info "â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”"
fi

# Auto-create API token for Home Assistant integration if enabled
if [[ "${auto_create_api_token}" == "true" && "${home_assistant_integration}" == "true" ]]; then
    bashio::log.info "Setting up Home Assistant integration..."
    ./bin/rails runner "
      # Create API token if it doesn't exist
      unless ApiToken.exists?(name: 'Home Assistant Addon')
        token = ApiToken.generate_for_integration('Home Assistant Addon')
        puts 'Created API token for Home Assistant integration'
        puts 'Token: ' + token.token
      end

      # Set integration settings
      Setting.set('home_assistant_enabled', true)
      Setting.set('api_enabled', true)
      Setting.set('webhook_base_url', ENV.fetch('WEBHOOK_BASE_URL', 'http://homeassistant.local:3000'))

      puts 'Home Assistant integration configured'
    " 2>/dev/null || bashio::log.warning "Failed to setup integration (this is normal on first run)"
fi

# Start the Rails server
bashio::log.info "Starting HomeChat server on port ${port}..."

if [[ "${ssl_enabled}" == "true" ]]; then
    bashio::log.info "Starting with SSL enabled..."
    # Start Puma directly with SSL configuration
    exec bundle exec puma -b "ssl://0.0.0.0:${port}?key=/ssl/${keyfile}&cert=/ssl/${certfile}" -e production
else
    exec ./bin/rails server -b 0.0.0.0 -p "${port}"
fi